# backend/scangrid.service
[Unit]
Description=ScanGRID API Server - Gridfinity Inventory Management
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/ScanGRID/backend
Environment="PATH=/home/pi/ScanGRID/backend/venv/bin"
Environment="SCANGRID_DB_DIR=/var/lib/scangrid"
ExecStart=/home/pi/ScanGRID/backend/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000 --log-level info
Restart=always
RestartSec=10

# Sécurité
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=scangrid

[Install]
WantedBy=multi-user.target
